# templates/hooks/deb/preinst.in
#!/bin/sh

@common_linux_functions@

cleanup_script_install() {
    echo "Cleaning up existing Sumo Otelcol installations (migration to native packaging) ..."

    # stop and disable systemd service
    if [ -f "/etc/systemd/system/otelcol-sumo.service" ]; then
        echo "Stopping and disabling existing Sumo Otelcol service ..."
        systemctl stop otelcol-sumo || true
        systemctl disable otelcol-sumo || true
        rm -f /etc/systemd/system/otelcol-sumo.service
    fi

    # remove binary
    if [ -f "/usr/local/bin/otelcol-sumo" ]; then
        echo "Removing existing Sumo Otelcol binary ..."
        rm -f /usr/local/bin/otelcol-sumo
    fi

    # remove configuration
    if [ -d "/etc/otelcol-sumo" ]; then
        echo "Removing existing Sumo Otelcol configuration directory (/etc/otelcol-sumo) ..."
        rm -rf /etc/otelcol-sumo
    fi

    # remove data
    if [ -d "/var/lib/otelcol-sumo" ]; then
        echo "Removing existing Sumo Otelcol data directory (/var/lib/otelcol-sumo) ..."
        rm -rf /var/lib/otelcol-sumo
    fi
}

case "$1" in
    install)
        cleanup_script_install
        create_group_if_missing
        create_user_if_missing
    ;;

    upgrade)
        create_group_if_missing
        create_user_if_missing
    ;;

    abort-upgrade)
    ;;

    *)
    echo "preinst called with unknown argument \`$1'" >&2
    exit 1
    ;;
esac

exit 0
